# scripts/systemd/rava-docs.service
[Unit]
Description=theexonet Game Docs
After=network.target

[Service]
WorkingDirectory=/var/www/publish
ExecStart=/usr/bin/dotnet /var/www/publish/Rava.Docs.dll
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=ASPNETCORE_URLS=http://0.0.0.0:9000
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_ENVIRONMENT=Production
Environment=RAVA_DATA_DIR=/var/www/data
User=www-data

[Install]
WantedBy=multi-user.target
