# scripts/systemd/rava-admin.service
[Unit]
Description=RAVA Admin Portal
After=network.target rava-api.service

[Service]
WorkingDirectory=/var/www/publish
ExecStart=/usr/bin/dotnet /var/www/publish/Rava.Admin.dll
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=ASPNETCORE_URLS=http://0.0.0.0:7000
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_ENVIRONMENT=Production
Environment=RAVA_DATA_DIR=/var/www/data
User=www-data

[Install]
WantedBy=multi-user.target
