# scripts/systemd/rava-permissions.service
[Unit]
Description=RAVA hosting permissions watcher
Documentation=file:///usr/local/lib/rava/scripts/fix-hosting-permissions.sh
After=network.target rava-api.service
Wants=rava-api.service

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/lib/rava/scripts/rava-permissions-watch.sh
EnvironmentFile=-/etc/default/rava-permissions
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rava-permissions

# Avoid restart storms if scripts are missing during install.
StartLimitIntervalSec=120
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
