# scripts/systemd/rava-admin.service
[Unit]
Description=RAVA Admin Portal
After=network.target rava-api.service

[Service]
WorkingDirectory=/var/www/publish
ExecStart=/usr/bin/dotnet /var/www/publish/Rava.Admin.dll
Restart=always
Environment=ASPNETCORE_URLS=http://0.0.0.0:7000
Environment=DOTNET_ENVIRONMENT=Production
User=www-data

[Install]
WantedBy=multi-user.target
