# scripts/systemd/theexonet-moderator.service
[Unit]
Description=theexonet Moderator Portal
After=network.target

[Service]
WorkingDirectory=/var/www/publish
ExecStart=/usr/bin/dotnet /var/www/publish/Theexonet.Moderator.dll
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=ASPNETCORE_URLS=http://0.0.0.0:7050
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_ENVIRONMENT=Production
Environment=THEEXONET_DATA_DIR=/var/www/data
User=theexonet
Group=theexonet

[Install]
WantedBy=multi-user.target
